# srcs/ftps/configs/supervisord.conf
[supervisord]
nodaemon = true
user=root

[program:ftps]
command=vsftpd
autorestart=true
startretries=3
user=root

[program:setup]
command=sh /dir/boot.sh
startsecs=0
startretries=0
user=root
