# srcs/grafana/srcs/supervisord.conf
[supervisord]
nodaemon = true
user=root

[program:grafana]
directory=/grafana/bin/
command=./grafana-server
autorestart=true
startretries=3
user=root
